# enigma2/var/spool/cron/crontabs/root
* * * * * if [ `cat /proc/stb/avs/0/input` = scart ]; then  date +%H%M > /dev/vfd; fi
* 23 * * 1 
